# deploy/systemd/erpnext-backup.timer
[Unit]
Description=Run daily backup for erpnext

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
